# etc/supervisord.d/n2n.conf
[program:n2n]
command=/usr/bin/supernode -l 54321
autorestart=true
stdout_logfile=/var/log/n2n.stdout.log
stderr_logfile=/var/log/n2n.stderr.log
